# systemd/update_pubkeys.timer
[Unit]
Description=Daily update of SSH pubkeys

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
